# restream/wc-restream.service
# /etc/systemd/system/wc-restream.service
# Keeps the restream running forever (boot + crash auto-restart).
#   scp restream.sh -> /opt/restream.sh ; chmod +x
#   scp wc-restream.service -> /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now wc-restream
# Change the channel by editing ExecStart args:  /opt/restream.sh 301 bein-max-2
[Unit]
Description=goattv -> HLS restream (ffmpeg)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/opt/restream.sh 299 bein-max-1
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
