# restream/wc-restream.service
# /etc/systemd/system/wc-restream.service
# Keeps the restream running forever (boot + crash auto-restart).
#   sudo cp restream.sh /opt/restream.sh && sudo chmod +x /opt/restream.sh
#   sudo cp wc-restream.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now wc-restream
# Change the channel by editing ExecStart args:  /opt/restream.sh 301 bein-max-2
[Unit]
Description=goattv -> HLS restream (ffmpeg)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/opt/restream.sh 301 bein-max-2
Restart=always
RestartSec=5
User=www-data
# ffmpeg + nginx write here
RuntimeDirectory=hls

[Install]
WantedBy=multi-user.target
